# docker/ops-git/rootfs/etc/services.d/70-tail/run
#!/usr/bin/with-contenv bash

touch /tmp/update.log
chmod 777 /tmp/update.log
exec tail -f /tmp/update.log