# docker/ops-git/rootfs/etc/services.d/70-tail/run
#!/usr/bin/with-contenv bash

# Wait for log file to be created by cron setup
sleep 2

touch /var/log/git-sync.log
exec tail -f /var/log/git-sync.log
